# web/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f -L 15

[program:php-fpm]
command=php-fpm -F

[inet_http_server]
port = *:9001
username = chenxiaojing
password = 8ik,9ol.
